# package/app/app/protocol-collect-uart/files/root/etc/uci-defaults/uci-defaults-protocol-collect-uart
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@protocol-collect-uart[-1]
	add ucitrack protocol-collect-uart
        set ucitrack.@protocol-collect-uart[-1].init=protocol-collect-uart
#        add_list ucitrack.@sub-ghz[0].init=sub-ghz
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
